# deploy/aemo-historical-prices.timer
[Unit]
Description=Monitor AEMO Historical Prices for updated monthly data
Documentation=https://github.com/cutout-z/aemo-historical-prices

[Timer]
# Check daily so new AEMO monthly files are picked up soon after publication.
OnCalendar=*-*-* 00:25:00
RandomizedDelaySec=20m
Persistent=true
Unit=aemo-historical-prices.service

[Install]
WantedBy=timers.target
